# var/www/vnc/audio_server/etc/init.d/vnc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
DISPLAYS="admin:2"

depend() {
	need net
}



start() {
	ebegin "Starting VNC server"

	for user in "admin:2"; do
		usrname=${user%%:*}
		usrdisp=${user##*:}
		su -s /bin/bash ${user%%:*} -c "source /etc/profile && cd ~${user%%:*} && export FD_GEOM=1920x1080 && x11vnc -shared -many -nomodtweak -noxkb -create -rfbport 590${usrdisp} -rfbauth ~/.vnc/passwd" >/var/log/vncserver.${usrdisp}.log 2>&1 &
                websockify --cert /etc/cert.pem 580${usrdisp} localhost:590${usrdisp} >/var/log/websockify.vncserver.${user%%:*}.log 2>&1 &
                pipeAudio ${user%%:*} 570${usrdisp}
	done
	eend $?
}

pipeAudio(){
 su $1 -c "/bin/tcpulse $2 /etc/cert.pem" >/var/log/tcpulse.$1.log 2>&1 &
}

stop() {
	ebegin "Stopping VNC server"
        killall websockify
        killall tcpulse
        pgrep gst-launch | xargs kill -9
        pkill -U 1000

	for user in $DISPLAYS; do
	     killall x11vnc >/dev/null 2>&1
	done
        rm -rf /tmp/.X11-unix >/dev/null 2>&1
 	eend 0
}

restart() {
        svc_stop
        svc_start
}
