# etc/init.d/vnc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
DISPLAYS="admin:2"

depend() {
	need net
}

start() {
	ebegin "Starting VNC server"

	for user in "admin:2"; do
		usrname=${user%%:*}
		usrdisp=${user##*:}
		su -s "/bin/bash" $usrname -c "startx -- :$usrdisp" > /var/log/x11.$usrname.log 2>&1 &
		su -s /bin/bash ${usrname} -c "source /etc/profile && export UNIXPW_DISABLE_SSL=1 && export UNIXPW_DISABLE_LOCALHOST=1 && cd ~$usrname && export FD_GEOM=1920x1080 && x11vnc -shared -many -nomodtweak -noxkb -create -rfbport 590${usrdisp} -unixpw_cmd /bin/verify_pw -loop -display \":$usrdisp\" " >/var/log/vncserver.${usrdisp}.log 2>&1 &
		websockify --cert /etc/cert.pem 580${usrdisp} localhost:590${usrdisp} >/var/log/websockify.vncserver.$usrname.log 2>&1 &
		su $1 -c "/bin/tcpulse 570${usrdisp} /etc/cert.pem" >/var/log/tcpulse.$usrname.log 2>&1 &
	done
	eend $?
}


stop() {
	ebegin "Stopping VNC server"
        killall websockify
        killall tcpulse
        pkill -U 1000

        rm -rf /tmp/.X11-unix >/dev/null 2>&1
 	eend 0
}

restart() {
        svc_stop
        svc_start
}
